# services.d/s6-rc.d/redis/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

redis-server --unixsocket /var/run/redis.sock --unixsocketperm 775