# etc/s6-overlay/s6-rc.d/run-app/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app && python3 -m tmdb_notifier
